# rootfs/etc/supervisord.conf
[supervisord]
user=root
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D -e

[program:gak]
command=/usr/bin/gak
